# profile.d/01-unset-base.sh
unset PROMPT_COMMAND
unset GIT_DISCOVERY_ACROSS_FILESYSTEM
