# profile.d/25-umask.sh
# Set the default umask for an interactive session
umask 0022
